# core-lib/src/config/generate/default.service
[Unit]
Description=Starts hyprshell daemon
After=graphical-session.target
Wants=graphical-session.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
ExecStartPre=/bin/sh -c '[ "$XDG_CURRENT_DESKTOP" = "Hyprland" ] || exit 0'
ExecStart={path} run {params}
Restart=always
{extra}

[Install]
WantedBy=graphical-session.target
